# eng/linux/Generic/usr/lib/systemd/user/opentabletdriver.service
[Unit]
Description=OpenTabletDriver Daemon
After=graphical.target

[Service]
ExecStart=/usr/share/OpenTabletDriver/OpenTabletDriver.Daemon
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
